# vps-trafficd/packaging/vps-trafficd.service
[Unit]
Description=VPS traffic JSON service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vps-trafficd --config /etc/vps-trafficd/config.toml
Restart=on-failure
RestartSec=5s
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/vps-trafficd

[Install]
WantedBy=multi-user.target
